# sync/scheduler_templates/archive_three_realm_sync_template.xml
<?xml version="1.0" encoding="UTF-16"?>
<!--
  Archive Three-Realm Sync Daily Scheduler Template (DO-SUPERVISOR-003)

  Drafted: 2026-04-30 (Day 132 morning, supervisor instance A self-implementation)
  Basis:
  - DO-SUPERVISOR-003 (archive three-realm sync, order #44 section 3 rebuttal 2)
  - DO-FACTORY-166 prerequisite (supervisor_report.ts data source)
  - B-003 Decision 1 Phase 1 prerequisite
  - ADR-009 section 6 three-realm automation integration

  Schedule: Daily 02:55 JST (5 minutes before supervisor dream-crystallize 03:00 JST)
  Action:   sync/sync_script/sync-archive-three-realm.ps1
  Output:   ProjectRX_HQ/wt_common/record-x-mirror/{supervisor,commander}/...

  Manual installation:
    schtasks /Create /TN "RecordX_ArchiveSync_ThreeRealm" `
      /XML "C:\RX_Dev\record-x-supervisor\sync\scheduler_templates\archive_three_realm_sync_template.xml"
-->
<Task version="1.4" xmlns="http://schemas.microsoft.com/windows/2004/02/mit/task">
  <RegistrationInfo>
    <Description>Daily archive three-realm sync (supervisor + commander archives -> factory mirror). Prerequisite for DO-FACTORY-166 (supervisor_report.ts). One-way mirror, mtime-based delta transfer, idempotent.</Description>
    <Author>RecordX Supervisor (instance A, 2026-04-30 Day 132)</Author>
    <URI>\RecordX\ArchiveSync_ThreeRealm</URI>
  </RegistrationInfo>
  <Triggers>
    <CalendarTrigger>
      <StartBoundary>2026-05-01T02:55:00+09:00</StartBoundary>
      <Enabled>true</Enabled>
      <ScheduleByDay>
        <DaysInterval>1</DaysInterval>
      </ScheduleByDay>
    </CalendarTrigger>
  </Triggers>
  <Principals>
    <Principal id="Author">
      <LogonType>InteractiveToken</LogonType>
      <RunLevel>LeastPrivilege</RunLevel>
    </Principal>
  </Principals>
  <Settings>
    <MultipleInstancesPolicy>IgnoreNew</MultipleInstancesPolicy>
    <DisallowStartIfOnBatteries>false</DisallowStartIfOnBatteries>
    <StopIfGoingOnBatteries>false</StopIfGoingOnBatteries>
    <AllowHardTerminate>true</AllowHardTerminate>
    <StartWhenAvailable>true</StartWhenAvailable>
    <RunOnlyIfNetworkAvailable>false</RunOnlyIfNetworkAvailable>
    <IdleSettings>
      <StopOnIdleEnd>false</StopOnIdleEnd>
      <RestartOnIdle>false</RestartOnIdle>
    </IdleSettings>
    <AllowStartOnDemand>true</AllowStartOnDemand>
    <Enabled>false</Enabled>
    <Hidden>false</Hidden>
    <RunOnlyIfIdle>false</RunOnlyIfIdle>
    <WakeToRun>false</WakeToRun>
    <ExecutionTimeLimit>PT15M</ExecutionTimeLimit>
    <Priority>7</Priority>
  </Settings>
  <Actions Context="Author">
    <Exec>
      <Command>powershell.exe</Command>
      <Arguments>-NoProfile -ExecutionPolicy Bypass -File "C:\RX_Dev\record-x-supervisor\sync\sync_script\sync-archive-three-realm.ps1"</Arguments>
      <WorkingDirectory>C:\RX_Dev\record-x-supervisor</WorkingDirectory>
    </Exec>
  </Actions>
</Task>
